# image-factory/images/bundle-trusty-chatwatt/ansible/files/lets-chat.conf
# lets-chat.conf

description "A chat app for small teams."

start on started networking
stop on runlevel [016]


limit nofile 1000000 1000000


console log

script

  mkdir -p /opt/lets-chat
  cd /opt/lets-chat

  # export NODE_PATH=/usr/lib/node_modules
  sudo npm start
end script

respawn
